# ops/ec2/finproof-agent-analysis-worker.service
[Unit]
Description=FinProof Agent analysis worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=finproof
Group=finproof
WorkingDirectory=/opt/finproof-agent/current
EnvironmentFile=/etc/finproof-agent/finproof-agent.env
Environment=NODE_ENV=production
Environment=NEXT_TELEMETRY_DISABLED=1
ExecStart=/usr/bin/env npm run ops:analysis:worker -- --loop
Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGINT
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
